# data/server/service_files/king-phisher.conf
# KingPhisherServer upstart config file
# This file goes in /etc/init/

description "King-Phisher Server"

start on runlevel [2345]
stop on runlevel [!2345]

exec /opt/king-phisher/KingPhisherServer -f /opt/king-phisher/server.conf
