# data/server/service_files/king-phisher.conf
# KingPhisherServer upstart config file
# This file goes in /etc/init/

description "King-Phisher Server"

start on runlevel [2345]
stop on runlevel [!2345]

# Installed Together In /opt/king-phisher
# exec /opt/king-phisher/KingPhisherServer -f /opt/king-phisher/server.conf
# Installed With 'python server-setup.py install'
exec KingPhisherServer -f /etc/king-phisher.conf
